# steam-flatpak-launcher-create.path
[Unit]
Description="Check the flatpak user data folder, and run the service accompanying this unit"

[Path]
PathModified=/home/steam/.local/share/flatpak/app
Unit=steam-launcher-create.service
